# deploy/mi-backup.timer
[Unit]
Description=Run mi state backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
